# scripts/crontab/preview
#
# !!AUTO-GENERATED!!  Edit scripts/crontab/make-crons.py instead.
#
MAILTO=flightdeck-developers@mozilla.org

#once per day
30 1 * * * cd /data/amo_python/src/builder.preview/flightdeck; /usr/bin/python26 manage.py cron clean_tmp
#every hour
30 * * * * cd /data/amo_python/src/builder.preview/flightdeck; /usr/bin/python26 manage.py cron celery

MAILTO=root
